# notebookcloud.conf
start on (stopped rc RUNLEVEL=[2345] and stopped cloud-config)

setuid ubuntu

exec /home/ubuntu/.notebookcloud_autofire.py
